# deploy/news-agent.timer
[Unit]
Description=Run News Agent twice daily at Toronto time

[Timer]
OnCalendar=*-*-* 06,18:00:00 America/Toronto
Persistent=true

[Install]
WantedBy=timers.target
