# deploy/news-agent.timer
[Unit]
Description=Run News Agent every 6 hours at Toronto time

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00 America/Toronto
Persistent=true

[Install]
WantedBy=timers.target
